# scheduler/crontab
PATH=/app/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
PYTHONUNBUFFERED=1
*/10 * * * * root /app/scheduler/run_once.sh >> /proc/1/fd/1 2>> /proc/1/fd/2

